# com.auto-web.generate.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.auto-web.generate</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/jcxs2014/Sites/Workbuddy/auto-web/gen_all.sh</string>
  </array>

  <!-- 每天本地时间 08:00 运行；如需其他时段改 Hour/Minute 即可 -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>8</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>WorkingDirectory</key>
  <string>/Users/jcxs2014/Sites/Workbuddy/auto-web</string>

  <!-- 日志便于排查（生成耗时较长属正常） -->
  <key>StandardOutPath</key>
  <string>/Users/jcxs2014/Sites/Workbuddy/auto-web/gen_all.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/jcxs2014/Sites/Workbuddy/auto-web/gen_all.err</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/Users/jcxs2014/.workbuddy/binaries/python/versions/3.13.12/bin</string>
  </dict>

  <!-- 机器唤醒/停电错过时间后补跑 -->
  <key>RunAtLoad</key>
  <true/>
</dict>
</plist>
